# root/crontab
* * * * * /bin/bash /root/elevatedWebControls.sh
